# infrastructure/ansible/roles/redbutton/files/systemd/xorg.service
[Unit]
Description=Launch X

[Service]
Type=simple
ExecStart=/usr/bin/X -nolisten tcp -novtswitch :0

[Install]
WantedBy=multi-user.target
